# app/crontab
#run script 1x per day to change guest new work password
0 0 * * * /usr/bin/python3 app/cron.py
#new line