# app/crontab
# Job for changing guest password on rotation
# @daily root /usr/local/bin/python3 /app/app/cron.py
@daily root /usr/bin/python3 /app/app/config/rotatePW.py 4 # test
