# app/crontab
#run script 1x per day to change guest network password
0 0 * * * root /usr/local/bin/python3 /app/app/cron.py >> /tmp/cron_temp 2>&1

