# app/crontab
# Job for changing guest password on rotation
0 0 * * * root /usr/local/bin/python3 /app/app/cron.py

